# work/old/BILLING_CREATE_COLLECTION/billing.createCollection.cron
# userlog remove 2017-10-23 yhlim
02 00 */3 * * sh /home/dreamsearch/scripts/billing.createCollection.sh
